# systemd/run-media-nas.automount
[Unit]
Description=Selfhosted samba share automount

[Automount]
Where=/run/media/nas
TimeoutIdleSec=60

[Install]
WantedBy=multi-user.target
